# dotfiles/configs/systemd/user/git-sync-obsidian.timer
[Unit]
Description=Run Git Sync for Obsidian Vault

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
